# server/scripts/systemd/task-agent-recommend.timer
[Unit]
Description=Run task-agent-recommend daily at 07:00 JST

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
